# services/aiml-service/app/api/v1/aiml/agent/installer/linux/aicompetency-agent.service
[Unit]
Description=AI Competency Agent - Local Python execution agent
Documentation=https://github.com/aptor/ai-competency-platform
After=network.target

[Service]
Type=simple
User=%i
ExecStart=/usr/bin/aicompetency-agent --host 127.0.0.1 --port 8889
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aicompetency-agent

# Security settings
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/aicompetency-agent

# Resource limits
LimitNOFILE=65536

# Working directory
WorkingDirectory=/var/lib/aicompetency-agent

# Environment
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=default.target
