# services/aiml-service/app/api/v1/aiml/agent/installer/macos/com.aicompetency.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!--
        LaunchAgent plist for AI Competency Agent
        
        This file configures the agent to run automatically when the user logs in.
        It will be installed to ~/Library/LaunchAgents/ during package installation.
        
        The agent runs as a background process and automatically restarts if it crashes.
    -->
    
    <!-- Service label (must match filename without .plist) -->
    <key>Label</key>
    <string>com.aicompetency.agent</string>
    
    <!-- Program to run -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/aicompetency-agent</string>
        <string>--host</string>
        <string>127.0.0.1</string>
        <string>--port</string>
        <string>8889</string>
    </array>
    
    <!-- Run at load (when user logs in) -->
    <key>RunAtLoad</key>
    <true/>
    
    <!-- Keep alive - restart if process exits -->
    <key>KeepAlive</key>
    <dict>
        <!-- Restart if process crashes -->
        <key>SuccessfulExit</key>
        <false/>
        <!-- Restart if process exits with non-zero code -->
        <key>Crashed</key>
        <true/>
    </dict>
    
    <!-- Standard output and error logging -->
    <key>StandardOutPath</key>
    <string>/usr/local/var/log/aicompetency-agent/agent.log</string>
    
    <key>StandardErrorPath</key>
    <string>/usr/local/var/log/aicompetency-agent/agent.error.log</string>
    
    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/usr/local/var/lib/aicompetency-agent</string>
    
    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>
    
    <!-- Process type - background daemon -->
    <key>ProcessType</key>
    <string>Background</string>
    
    <!-- Throttle interval - wait 10 seconds before restarting -->
    <key>ThrottleInterval</key>
    <integer>10</integer>
    
    <!-- Disable automatic restart if it crashes too many times -->
    <key>StartInterval</key>
    <integer>10</integer>
    
    <!-- Limit restart attempts (optional - comment out for unlimited restarts) -->
    <!--
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>0</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    -->
</dict>
</plist>
